# server/chat/python/docs/systemd-with-workers/system/chat-server.target
[Unit]
Description=Chat server parent target
After=network.target

[Install]
WantedBy=multi-user.target
